# Python/3.6.8/Django-uWSGI-Nginx/Config/Supervisor.conf
[program:uwsgi]
command = /usr/local/bin/uwsgi --ini /home/docker/Project/uWSGI-Emperor.ini

[program:nginx]
command = /usr/sbin/nginx
